# scripts/systemd/excalidraw.service
[Unit]
Description=Excalidraw

[Service]
Type=simple
ExecStart=/usr/local/bin/excalidraw
Restart=always

[Install]
WantedBy=multi-user.target
